# packaging/macos/dev.union.client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.union.client</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Union.app/Contents/MacOS/union-client</string>
        <string>--config</string>
        <string>__HOME__/Library/Application Support/Union/client.toml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/Union/client.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/Union/client.err.log</string>
</dict>
</plist>
